# deploy/fas_smart.service
[Unit]
Description=FAS Smart Real-Time Signal Service
After=network.target postgresql.service

[Service]
Type=simple
User=elcrypto
WorkingDirectory=/home/elcrypto/fas_smart
Environment=PATH=/home/elcrypto/fas_smart/venv/bin
ExecStart=/home/elcrypto/fas_smart/venv/bin/python scripts/run_service.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fas_smart

# Resource limits
MemoryMax=4G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
